# files/71-10g-ethernet.rules
SUBSYSTEM=="net", ACTION=="add", DRIVERS=="ixgbe", KERNEL=="eth*", RUN+="/sbin/ethtool -G %k rx 4096"
